# deploy/oracle/plagiarism-compose.service
[Unit]
Description=Plagiarism Checker (Docker Compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/opt/plagiarism-checker/deploy/oracle
ExecStart=/usr/bin/docker compose -f docker-compose.oracle.yml up -d --build
ExecStop=/usr/bin/docker compose -f docker-compose.oracle.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
